# deploy/gcp/ansible/roles/haproxy/templates/consul-template.service.j2
[Unit]
Description=consul-template — HAProxy config renderer ({{ routing_combo }})
Documentation=https://github.com/hashicorp/consul-template
After=network-online.target haproxy.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/consul-template -config=/etc/consul-template/consul-template.hcl
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
